# luffy-launcher/debian/postrm
#!/bin/sh
set -e

# Stop and disable service
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    systemctl stop luffy-gateway.service || true
    systemctl disable luffy-gateway.service || true
    systemctl daemon-reload
fi

#DEBHELPER#
